# layers/meta-nvidia-container-block/recipes-support/balena-nvidia-cdi-setup/files/balena-nvidia-cdi-setup.service
[Unit]
Description=Generate nvidia CDI spec + enable balena-engine CDI feature
DefaultDependencies=no
# /etc/docker is a bind-mount from the data partition wired up by
# bind-etc-docker.service; without ordering after that unit, our script
# writes daemon.json into a read-only overlay and fails with EROFS.
Requires=bind-etc-docker.service
After=local-fs.target bind-etc-docker.service
Before=balena.service balena-engine.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/libexec/balena-nvidia-cdi-setup

[Install]
WantedBy=balena.service
