# layers/meta-balena-jetson/recipes-support/os-fan-profile/os-fan-profile/etc-nvfancontrol-config.mount
[Unit]
Description=Bind mount for /etc/nvfancontrol.conf
Requires=resin-state.service resin-state-reset.service
After=resin-state.service resin-state-reset.service

[Mount]
What=/mnt/state/root-overlay/etc/nvfancontrol-config
Where=/etc/nvfancontrol-config
Type=none
Options=bind

[Install]
WantedBy=multi-user.target
